# rpi-timelapsephoto/octopuslab-timelapsephoto.timer
[Unit]
Description=OctopusLAB Photoservice timer

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
